# deploy/uat/supervisord-uat-service.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:casvardb-api]
directory=/app/backend-nestjs
command=/bin/bash -lc 'PORT="${API_PORT:-8888}" npm start'
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:casvardb-worker]
directory=/app/backend-nestjs
command=/bin/bash -lc 'npm run start:worker'
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:casvardb-fe]
directory=/app/frontend
command=/bin/bash -lc 'HOST=0.0.0.0 PORT="${FRONTEND_PORT:-3000}" npm start'
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
